# upstart-configs/sdtd-tornado.conf
description "7D2D Tornado Site"
author "RaTilicus"
  
start on runlevel [2345]
stop on runlevel [06]
  
respawn
exec /var/www/sdtd/sdtd-tornado.py

