# upstart-configs/sdtd-log.conf
description "SDTD Telnet Logging"
author "RaTilicus"

setuid www-data
setgid www-data

start on runlevel [2345]
stop on runlevel [06]
  
respawn
exec /var/www/sdtd/sdtd-log.py

